# net-misc/smcroute/files/smcroute.initd
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
}

start() {
	ebegin "Starting smcroute"
	/usr/sbin/smcroute -d
	eend $?
}

stop() {
	ebegin "Stopping smcroute"
	/usr/sbin/smcroute -k
	eend $?

}
